# cookbooks/dude/templates/default/dude.conf.erb
# /etc/init/dude.conf
# THIS FILE IS MANAGED BY CHEF

# Once deployed, you can manage with:
# sudo start dude
# sudo stop dude
# sudo restart dude
# sudo status dude

description "Dude Daemon"

setuid <%= @user %>
setgid <%= @group %>
env HOME=/home/<%= @user %>
env ENV=production

start on runlevel [23]
stop on runlevel [06]

script
# should run in bash for chruby support
exec /bin/bash <<EOF
  . /etc/profile.d/chruby.sh
  cd <%= @cwd %>
  exec bundle exec lita
EOF
end script
